# packaging/scripts/after-install
#!/bin/sh

chown -R kotori:kotori /opt/kotori

systemctl reenable kotori
systemctl restart kotori
