# packaging/scripts/after-install
#!/bin/sh

systemctl reenable kotori
systemctl restart kotori
