# packaging/scripts/before-remove
#!/bin/sh

systemctl stop kotori
systemctl disable kotori
